# script/init/boxque.conf
description "Boxque worker"

start on runlevel [23]
stop on shutdown

exec /usr/local/bin/boxqueWorker

respawn